# rpm/80-hotplug-cpu-mem.rules
# do not edit this file, it will be overwritten on update

# Hotplug physical CPU
SUBSYSTEM=="cpu", ACTION=="add", TEST=="online", ATTR{online}=="0", \
 ATTR{online}="1"

# Hotplug physical memory. Make sure that the memory is onlined as movable
# This is only applicable to setups where all the onlined memory is aimed
# to be hotplugable (aka movable).
SUBSYSTEM=="memory", ACTION=="add", TEST=="state", ATTR{state}=="offline", \
 ATTR{state}="online_movable", TAG+="tmpfs"

# See bnc#869603
TAG=="tmpfs", RUN+="/usr/lib/udev/remount-tmpfs"
